# deploy/weather-app.service.example
[Unit]
Description=Weather app
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/micu/weather/anotherGo_at_weatherApp_js
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8105
Environment=TRUST_PROXY=true
Environment=APP_VERSION=1.0.0
EnvironmentFile=/home/micu/weather/anotherGo_at_weatherApp_js/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
